# data/LaunchDaemons/jp.imoz.imos-start.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>KeepAlive</key>
	<false/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>IMOSH_FLAGS_alsologtostderr</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>jp.imoz.imos-start</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/usr/local/imos/bin/imos-start</string>
		<string>--alsologtostderr</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/tmp/jp.imoz.imos-start.log.STDOUT</string>
	<key>StandardErrorPath</key>
	<string>/tmp/jp.imoz.imos-start.log.STDERR</string>
</dict>
</plist>
